# templates/systemd_timer.timer.j2
[Unit]
Description={{ description }}

[Timer]
OnCalendar={{ on_calendar }}
Persistent={{ persistent }}

[Install]
WantedBy={{ wanted_by }}
